# config/daily-advisor.timer
[Unit]
Description=Run Claude Ops Daily Advisor

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
