# config/health-auditor.timer
[Unit]
Description=Run Claude HUD Health Auditor daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
